# salt/roots/salt/apps/datal/supervisor/uwsgi.conf
[program:uwsgi]
command={{ pillar['virtualenv']['path'] }}/bin/uwsgi --uid {{ pillar['system']['user'] }} --home {{ pillar['virtualenv']['path'] }} --emperor /etc/uwsgi/ --emperor-stats 127.0.0.1:3114 --logto /var/log/uwsgi/uwsgi-emperor.log
autostart=true
autorestart=true
stderr_logfile = /var/log/uwsgi/uwsgi-err.log
stdout_logfile = /var/log/uwsgi/uwsgi-out.log